# packages/hnefatafl-ai-attacker.service
[Unit]
Description=Copenhagen Hnefatafl AI
Documentation="https://hnefatafl.org"

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/hnefatafl-ai --username ${USERNAME} --password ${PASSWORD} --role attacker
PrivateTmp=yes
InaccessibleDirectories=/home
DeviceAllow=/dev/null rw
EnvironmentFile=/etc/hnefatafl-ai-attacker.conf

[Install]
WantedBy=multi-user.target
