# packages/hnefatafl-ai-attacker.service
[Unit]
Description=Copenhagen Hnefatafl AI
Documentation="https://hnefatafl.org"

[Service]
Type=simple
Restart=always
ExecStart=/usr/games/taflzero --username ${USERNAME} --password ${PASSWORD} --role attacker --accept-game 1 --systemd
PrivateTmp=yes
InaccessibleDirectories=/home
DeviceAllow=/dev/null rw
EnvironmentFile=/etc/hnefatafl-ai-attacker.conf
WorkingDirectory=/usr/share/taflzero

[Install]
WantedBy=multi-user.target
