# packages/hnefatafl.service
[Unit]
Description=Copenhagen Hnefatafl Server
Documentation="https://hnefatafl.org"

[Service]
Type=simple
Restart=always
ExecStart=/usr/games/hnefatafl-server-full --systemd --secure --autostart-tournament
PrivateTmp=yes
InaccessibleDirectories=/home
DeviceAllow=/dev/null rw

[Install]
WantedBy=multi-user.target
