# packages/hnefatafl.service
[Unit]
Description=Copenhagen Hnefatafl Server
Documentation="https://hnefatafl.org"

[Service]
Type=simple
Restart=always
ExecStart=/usr/games/hnefatafl-server-full --systemd

[Install]
WantedBy=multi-user.target
