# packages/hnefatafl-ai-00-attacker.service
[Unit]
Description=Copenhagen Hnefatafl AI
Documentation="https://hnefatafl.org"

[Service]
Type=simple
Restart=always
ExecStart=/usr/games/hnefatafl-ai --username 00 --password PASSWORD --role attacker
PrivateTmp=yes
InaccessibleDirectories=/home
DeviceAllow=/dev/null rw

[Install]
WantedBy=multi-user.target
