# packages/hnefatafl-ai-defender.service
[Unit]
Description=Copenhagen Hnefatafl AI
Documentation="https://hnefatafl.org"

[Service]
Type=simple
Restart=always
ExecStart=/usr/games/hnefatafl-ai --username ${USERNAME}  --password ${PASSWORD} --role defender --parallel
PrivateTmp=yes
InaccessibleDirectories=/home
DeviceAllow=/dev/null rw
EnvironmentFile=/etc/hnefatafl-ai-defender.conf

[Install]
WantedBy=multi-user.target
